# dist/signd.service
[Unit]
Description=GPG Sign Daemon
After=syslog.target obsapisetup.service

[Service]
Type=forking
PIDFile=/run/signd.pid
EnvironmentFile=-/etc/sysconfig/signd
ExecStart=/usr/sbin/signd -f
Restart=on-abort

[Install]
WantedBy=multi-user.target
